# nginx/supervisor/conf.d/puppet.conf
[program:puppet]
command=/scripts/puppet.sh
numprocs=1
autostart=true
autorestart=true

