# puppetmaster/supervisor/conf.d/puppetmaster.conf
[program:puppetmaster]
command=/scripts/puppetmaster.sh
numprocs=1
autostart=true
autorestart=true

