# dnsmasq/supervisor/conf.d/dnsmasq.conf
[program:dnsmasq]
command=/scripts/dnsmasq.sh
numprocs=1
autostart=true
autorestart=true

