# dnsmasq/supervisor/conf.d/serf.conf
[program:serf]
command=/scripts/serf.sh
numprocs=1
autostart=true
autorestart=true

